# molecule/python3/files/pycsw.conf
[program:pycsw-load]
command={{ python_app_virtualenv_bin }}/pycsw-load.sh
autostart=false
autorestart=false
