# molecule/default/files/supervisor.conf
[program:app]
command=echo hi
autostart=false
autorestart=false
